# supervisor-godns.conf
[program:godns]
command =/app/godns -c /etc/godns.conf

[inet_http_server]
port = :9001