# supervisor-joke.conf
[program:joke]
command = /app/joke -c /etc/joke.conf

[inet_http_server]
port = :9001